# ops/space.hackertrip.review-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>space.hackertrip.review-sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/npx</string>
    <string>tsx</string>
    <string>scripts/sync-admin-to-feishu.ts</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/jaydenworkplace/hackertrip-recovery</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/jaydenworkplace/Library/Logs/hackertrip-review-sync.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/jaydenworkplace/Library/Logs/hackertrip-review-sync.log</string>
</dict>
</plist>
