# ops/space.hackertrip.daily-ops.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>space.hackertrip.daily-ops</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/npx</string>
    <string>tsx</string>
    <string>scripts/run-daily-hackathon-ops.ts</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/jaydenworkplace/Desktop/hacker_trip</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/jaydenworkplace/Library/Logs/hackertrip-daily-ops.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/jaydenworkplace/Library/Logs/hackertrip-daily-ops.log</string>
</dict>
</plist>
